# makeshinylinks.conf
# This is the upstart script to make links to shiny applications

description "Shiny application mapping"

start on runlevel [2345]
stop on runlevel [016]

exec /vagrant/makeshinylinks.sh

